# scripts/refine-runner.service
# Example systemd unit for refine-runner. Install at:
#   /etc/systemd/system/refine-runner.service
#
# Adjust:
#   - User/Group         — operator account (the one that ran `claude login`
#                          and has rights to the bound client repo)
#   - WorkingDirectory   — the refine source clone (the dir containing
#                          .refine-binding produced by `refine init`)
#
# refine has zero env-var configuration: the binding inside WorkingDirectory
# resolves which client to operate on.

[Unit]
Description=refine host runner
After=network.target

[Service]
Type=simple
User=refine
Group=refine
WorkingDirectory=/opt/refine-acme
ExecStart=/usr/local/bin/uv run refine runner --foreground
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
