# scripts/refine-runner.service
# Example systemd unit for refine-runner. Install at:
#   /etc/systemd/system/refine-runner.service
#
# Edit:
#   - User/Group           — operator account (must have run `claude login`)
#   - WorkingDirectory     — the client repo root (where ./refine/refine.toml lives)
#   - ExecStart's PYTHONPATH (or installed location of refine)
#
# refine has zero env-var configuration: refine.toml in WorkingDirectory does it.

[Unit]
Description=refine host runner
After=network.target

[Service]
Type=simple
User=refine
Group=refine
WorkingDirectory=/srv/clients/acme-app
# Adjust PYTHONPATH to point at where this codebase is installed on the host.
Environment=PYTHONPATH=/opt/refine
ExecStart=/usr/bin/python3 -m refine runner
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
